# .config/systemd/user/lspmux.socket
[Unit]
Description=lspmux socket
# StopWhenUnneeded=yes

[Socket]
ListenStream=%t/lspmux.sock
SocketMode=0600

# [Install]
# WantedBy=sockets.target
